# systemd/kfxai-blog.timer
[Unit]
Description=Run kfxai blog post at 05:30 / 13:30 / 21:30 JST (kfreqaiгБЃ05/13/21:00гБЛгВЙзіД30еИЖгБЪгВЙгБЧ)

[Timer]
OnCalendar=*-*-* 05,13,21:30:00
RandomizedDelaySec=180
Persistent=true

[Install]
WantedBy=timers.target
